# scripts/rpi/91-btremove.rules
ACTION=="remove", KERNEL=="hci0", RUN+="/usr/local/bin/udevstopbt.sh"
