# scripts/pc/90-bnep-start.rules
ACTION=="add", KERNEL=="bnep0", RUN+="/sbin/dhclient -v bnep0"
